# workingos/aios/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:xvfb]
command=Xvfb :0 -screen 0 1280x800x24 -ac -nolisten tcp +extension RANDR
priority=1
autorestart=true
stdout_logfile=/var/log/supervisor/xvfb.log
stderr_logfile=/var/log/supervisor/xvfb.err

[program:desktop]
command=dbus-run-session startxfce4
environment=DISPLAY=:0,HOME=/home/agent,USER=agent
user=agent
priority=2
autorestart=true
stdout_logfile=/var/log/supervisor/desktop.log
stderr_logfile=/var/log/supervisor/desktop.err

[program:vnc]
command=x11vnc -display :0 -forever -nopw -rfbport 5900 -shared -repeat
user=agent
priority=3
autorestart=true
stdout_logfile=/var/log/supervisor/vnc.log
stderr_logfile=/var/log/supervisor/vnc.err

[program:novnc]
command=websockify --web=/usr/share/novnc/ 6080 localhost:5900
priority=4
autorestart=true
stdout_logfile=/var/log/supervisor/novnc.log
stderr_logfile=/var/log/supervisor/novnc.err
